# supervisord/conf.d/app.conf
[program:app]
command=/usr/sbin/uwsgi -s /tmp/app.sock --manage-script-name --mount /=vhost_api:app
directory=/src/vhost_api
autostart=true
autorestart=true
priority=5
redirect_stderr=true
